# taco/misc/taco-manager.sh.in
#!/bin/sh
# Toolkit for building distributed control systems or any other distributed system.
#
# Copyright (c) 2005-2013 Jörg Pulz <joerg.pulz@frm2.tum.de>
#
# This library is free software; you can redistribute it and/or
# modify it under the terms of the GNU Lesser General Public
# License as published by the Free Software Foundation; either
# version 2.1 of the License, or (at your option) any later version.
#
# This library is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the GNU
# Lesser General Public License for more details.
#
# You should have received a copy of the GNU Lesser General Public
# License along with this library; if not, write to the Free Software
# Foundation, Inc., 51 Franklin St, Fifth Floor, Boston, MA  02110-1301  USA
#
# This is the FreeBSD specific init script
# 
# PROVIDE: taco_manager
# REQUIRE: DAEMON rpcbind
# BEFORE: LOGIN
# KEYWORD: FreeBSD shutdown
#
# Add the following lines to /etc/rc.conf to enable taco-manager:
#
# taco_manager_enable="YES"
#

. /etc/rc.subr

name=taco_manager
rcvar=`set_rcvar`

prefix=@prefix@
exec_prefix=@exec_prefix@
command=@sbindir@/Manager
required_files=@sysconfdir@/tacoenv.sh
start_precmd="taco_precmd"
start_postcmd="taco_postcmd"

taco_precmd()
{
	. @sysconfdir@/tacoenv.sh
	command_args=
	if test "x$TACO_LOGGING" = "xyes" ; then
		command_args="-log"
		if test ! -d $LOGPATH ; then
			mkdir -p $LOGPATH
		fi
        fi
	if test "x$TACO_SECURITY" = "xyes" ; then
		command_args="$command_args -security"
	fi
	command_args="$command_args -${TACO_DATABASE:-dbm}"
}

taco_postcmd()
{
	. @sysconfdir@/tacoenv.sh
	if [ ${TACO_DATABASE:-dbm} = dbm ] ; then
		new=`ls -l ${DBM_DIR}/* 2>/dev/null | awk '{print $5;}' | sort | uniq | wc -l`
		if [ $new -le 1 ] ;  then
			info "Initialization of the database:"
			if db_fillup2 ; then
				info "TACO database successfully loaded"
				return 0
			else
				err "TACO database not completely loaded"
				return 1
			fi
		fi
	elif test x"${TACO_DATABASE}" = x"sqlite3" ; then
		new=`sqlite3 ${DBM_DIR}/tango 'select count(*) from property_device'`
		if [ $new -eq 0 ] ; then
			info "Initialization of the database:"
			if db_fillup2 ; then
				info "TACO database successfully loaded"
			else
				err "TACO database not completely loaded"
				return 1
			fi
		fi
	fi
}

# set defaults

taco_manager_enable=${taco_manager_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
