# taco/misc/taco-manager.sh.in
#!/bin/sh
#

# PROVIDE: taco_manager
# REQUIRE: DAEMON rpcbind
# BEFORE: LOGIN
# KEYWORD: FreeBSD shutdown

#
# Add the following lines to /etc/rc.conf to enable taco-manager:
#
#taco_manager_enable="YES"
#

. /etc/rc.subr

name=taco_manager
rcvar=`set_rcvar`

command=@sbindir@/Manager
required_files=@sysconfdir@/tacoenv.sh
start_precmd="taco_precmd"

taco_precmd()
{
	. @sysconfdir@/tacoenv.sh
}

# set defaults

taco_manager_enable=${taco_manager_enable:-"NO"}
taco_manager_flags=${taco_manager_flags:-"-dbm &"}

load_rc_config $name
run_rc_command "$1"
