# taco/misc/taco-manager.sh.in
#!/bin/sh
#

# PROVIDE: taco_manager
# REQUIRE: DAEMON rpcbind
# BEFORE: LOGIN
# KEYWORD: FreeBSD shutdown

#
# Add the following lines to /etc/rc.conf to enable taco-manager:
#
#taco_manager_enable="YES"
#

. /etc/rc.subr

name=taco_manager
rcvar=`set_rcvar`

prefix=@prefix@
exec_prefix=@exec_prefix@
command=@sbindir@/Manager
required_files=@sysconfdir@/tacoenv.sh
start_precmd="taco_precmd"
start_post_cmd="taco_postcmd"

taco_precmd()
{
	. @sysconfdir@/tacoenv.sh
	taco_manager_flags=
	if test "x$TACO_LOGGING" = "xyes" ; then
		taco_manager_flags="-log"
		if ! test -d $LOGPATH ; then
			mkdir -p $LOGPATH
		fi
        fi
	if test "x$TACO_SECURITY" = "xyes" ; then
		taco_manager_flags="$taco_manager_flags -security"
	fi
	taco_manager_flags="$taco_manager_flags -${TACO_DATABASE:-dbm}"
}

taco_postcmd()
{
	sleep 2
}

# set defaults

taco_manager_enable=${taco_manager_enable:-"NO"}
# taco_manager_flags=${taco_manager_flags:-"-dbm"}

load_rc_config $name
run_rc_command "$1"
