# taco/misc/taco.debian.in
#! /bin/sh
# This script is Debian specific
### BEGIN INIT INFO
# Provides: taco
# Required-Start: $network $portmap
# Required-Stop:
# Default-Start: 3 5
# Default-Stop: 0 1 2 6
# Description: Start the TACO system services.
### END INIT INFO

usage="Usage: $0 {start|stop|restart|new}"

# Include TACO environment file
prefix=${DSHOME:-@prefix@}
exec_prefix=@exec_prefix@
taco_env_file="@sysconfdir@/tacoenv.sh"
if test ! -r "$taco_env_file" ; then
	echo "error: TACO environment file '$taco_env_file' is missing"
	exit 6
fi
. "$taco_env_file" 

# Check for TACO manager executable
taco_manager="${TACO_PATH}/Manager"
if test ! -x $taco_manager ; then
	echo "error: TACO manager '$taco_manager' is not installed"
	exit 5
fi

case $1 in
	start)
		if test "x$TACO_LOGGING" = xyes ; then 
			log="-log"
		fi
		if test "x$TACO_SECURITY" = xyes ; then
			sec_text="in secure mode"
			security="-security"
		fi
		dbase="-${TACO_DATABASE:-dbm}"
		echo -n "Starting TACO services $sec_text"
		start-stop-daemon --start --quiet --exec "$taco_manager -- $dbase $log $security"
		echo "."
		if [ ${TACO_DATABASE:-dbm} = dbm ] ; then
			new=`ls -l ${DBM_DIR}/* 2>/dev/null | awk '{print $5;}' | sort | uniq | wc -l`
			if [ $new -le 1 ] ;  then
				echo -n "Initialization of the database:"
				sleep 5
				if db_fillup2 ; then
					echo "."
				else
					exit 1
				fi
			fi
		fi
		;;
	stop)
		echo -n "Shutting down TACO services"
		start-stop-daemon --stop --quiet --oknodo --exec "$taco_manager"
		echo "."
		;;
	restart)
		$0 stop
		$0 start
		;;
	status)
		echo -n "Checking for TACO services"
		exit 3
		;;
	new)
		if [ ${TACO_DATABASE:-dbm} = dbm ] ; then
			echo -n "Do you really want to delete the complete TACO database (not the resource files)? (yes/no): "
			read
			if test "x$REPLY" = xyes ; then
				rm -f "$DBM_DIR"/*
				$0 restart
			fi
		fi
		;;
	*)
		echo "$usage"
		exit 1
		;;
esac
exit 0
