# taco/misc/taco.in
#! /bin/sh

### BEGIN INIT INFO
# Provides: taco
# Required-Start: $network $portmap
# Required-Stop:
# Default-Start: 3 5
# Default-Stop: 0 1 2 6
# Description: Start the TACO system services.
### END INIT INFO

. /etc/rc.status
rc_reset

usage="Usage: $0 {start|stop|status|restart|new}"

# Include TACO environment file
taco_env_file="/etc/tacoenv.sh"
if ! test -r "$taco_env_file" ; then
	echo "error: TACO environment file '$taco_env_file' is missing"
	exit 6
fi
. "$taco_env_file" 

# Check for TACO manager executable
taco_manager="${TACO_PATH}/Manager"
if ! test -x $taco_manager ; then
	echo "error: TACO manager '$taco_manager' is not installed"
	exit 5
fi

case $1 in
	start)
		if test ! \( -f "$DBM_DIR/names.pag" -o -f "$DBM_DIR/names" \) ; then
			echo "Initialization of the database:"
			if db_fillup 0 ; then
				rc_status -v
			else
				rc_status -v
				rc_exit
			fi
		fi
		if test "x$TACO_LOGGING" = xyes ; then 
			log="-log"
		fi
		if test "x$TACO_SECURE_MODE" = xyes ; then
			sec_text="in secure mode"
			security="-security"
		fi
		echo -n "Starting TACO services $sec_text"
		startproc -s -q "$taco_manager" $log $security
		rc_status -v
		;;
	stop)
		echo -n "Shutting down TACO services"
		killproc "$taco_manager"
		rc_status -v
		;;
	restart)
		$0 stop
		$0 start
		rc_status
		;;
	status)
		echo -n "Checking for TACO services"
		checkproc "$taco_manager"
		rc_status -v
		;;
	new)
		echo -n "Do you really want to delete the complete TACO database (not the resource files)? (yes/no): "
		read
		if test "x$REPLY" = xyes ; then
			rm -i "$DBM_DIR"/*
			$0 restart
			rc_status
		else
			rc_failed
		fi
		;;
	*)
		echo "$usage"
		rc_failed
		;;
esac
rc_exit
